# deploy/systemd/goyoonjung-wiki-daily.timer
[Unit]
Description=Run goyoonjung-wiki daily automation every 6 hours

[Timer]
OnBootSec=5min
OnUnitActiveSec=6h
AccuracySec=1min
Persistent=true
Unit=goyoonjung-wiki-daily.service

[Install]
WantedBy=timers.target
