# deploy/systemd/user/goyoonjung-wiki-watchdog.service
[Unit]
Description=goyoonjung-wiki watchdog automation
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/바탕화면/goyoonjung-wiki
ExecStart=/usr/bin/env bash -lc '%h/바탕화면/goyoonjung-wiki/scripts/watchdog_automation.py'
